# flask-app/transcriber.service
[Unit]
Description=Transcriber - AI-Powered Audio Transcription Service
After=network.target

[Service]
Type=simple
User=<user>
WorkingDirectory=/path/to/transcriber/flask-app
Environment="PATH=/path/to/transcriber/flask-app/.venv/bin"
EnvironmentFile=/path/to/transcriber/flask-app/.env
ExecStart=/path/to/transcriber/flask-app/.venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 300 --graceful-timeout 30 app:app
Restart=always
RestartSec=10
StandardOutput=append:/var/log/transcriber/transcriber.log
StandardError=append:/var/log/transcriber/transcriber-error.log

[Install]
WantedBy=multi-user.target
